# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

# Wireless radio0 specific settings
set wireless.radio0.htmode='HE40'
set wireless.radio0.channel='auto'
set wireless.radio0.country='PT'

# Wireless radio1 specific settings
set wireless.radio1.htmode='HE160'
set wireless.radio1.channel='auto'
set wireless.radio1.country='PT'

# Firewall settings
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'

# uHTTPd settings
set uhttpd.main.redirect_https='1'

EOI
